# installBrightness_SysFSEventListener.service
[Unit]
Description=Listen to Brightness events.

[Service]
ExecStart=/usr/local/bin/intelBrightness_HookSysFS

[Install]
WantedBy=sysinit.target
